# src/changeIP.cron
#!/bin/bash
cd /etc/bind
/home/styles/Installed/jython2.5.2/jython -Dpython.path="/home/styles/Installed/htmlunit-2.8/lib/httpclient-4.0.1.jar:/home/styles/Installed/htmlunit-2.8/lib/commons-collections-3.2.1.jar:/home/styles/Installed/htmlunit-2.8/lib/cssparser-0.9.5.jar:/home/styles/Installed/htmlunit-2.8/lib/httpcore-4.0.1.jar:/home/styles/Installed/htmlunit-2.8/lib/httpmime-4.0.1.jar:/home/styles/Installed/htmlunit-2.8/lib/apache-mime4j-0.6.jar:/home/styles/Installed/htmlunit-2.8/lib/commons-io-1.4.jar:/home/styles/Installed/htmlunit-2.8/lib/sac-1.3.jar:/home/styles/Installed/htmlunit-2.8/lib/nekohtml-1.9.14.jar:/home/styles/Installed/htmlunit-2.8/lib/xalan-2.7.1.jar:/home/styles/Installed/htmlunit-2.8/lib/serializer-2.7.1.jar:/home/styles/Installed/htmlunit-2.8/lib/commons-logging-1.1.1.jar:/home/styles/Installed/htmlunit-2.8/lib/xml-apis-1.3.04.jar:/home/styles/Installed/htmlunit-2.8/lib/commons-codec-1.4.jar:/home/styles/Installed/htmlunit-2.8/lib/xercesImpl-2.9.1.jar:/home/styles/Installed/htmlunit-2.8/lib/htmlunit-2.8.jar:/home/styles/Installed/htmlunit-2.8/lib/commons-lang-2.4.jar:/home/styles/Installed/htmlunit-2.8/lib/htmlunit-core-js-2.8.jar:/usr/share/pyshared" /etc/bind/changeIP.py

